# etc/udev/rules.d/75-waggle-net-interface.rules
SUBSYSTEM=="net", KERNELS=="6-1:1.0", RUN+="/bin/sed -i -e s/INTERFACE/$env{INTERFACE}/ /etc/network/interfaces"
