# etc/systemd/system/waggle-audio-collector.service
[Unit]
Description=Audio clip collector
After=waggle-platform.target

[Service]
ExecStart=/usr/bin/python3 -u /usr/lib/waggle/edge_processor/audio/pipeline/example_audio.py

[Install]
WantedBy=waggle-platform.target
