# etc/systemd/system/waggle-image-collector.service
[Unit]
Description=Image Collector

ConditionPathExists=/root/fs_locked

[Service]
ExecStart=/usr/bin/python3 /usr/lib/waggle/edge_processor/image/pipeline/image_collector.py

Restart=always
RestartSec=30